# system/crowker.service
[Unit]
Description=Crow brocker

[Service]
Type=simple
ExecStart=/usr/local/bin/crowker --udp 10009

[Install]
WantedBy=multi-user.target